# cfg/session-example/netidx.service
[Unit]
Description=Netidx Activation

[Service]
ExecStart=/home/eric/.cargo/bin/netidx activation -f

[Install]
WantedBy=default.target
